# build_assets/deb_files/content/etc/init/fiona.conf
description    fiona

start on runlevel [2345]
stop on runlevel [016]

chdir /opt/fiona
respawn
setuid fiona
setgid fiona
umask 0022

exec /opt/fiona/fiona.sh

